# scripts/local.kukanilea.ai-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.kukanilea.ai-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gensuminguyen/Tophandwerk/kukanilea-git/.venv312/bin/python3</string>
        <string>/Users/gensuminguyen/Tophandwerk/kukanilea-git/app/ai/sync_bridge.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/gensuminguyen/Tophandwerk/kukanilea-git</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer> <!-- Führt den Sync alle 60 Minuten aus -->
    <key>StandardOutPath</key>
    <string>/Users/gensuminguyen/Tophandwerk/kukanilea-git/reports/ai_sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gensuminguyen/Tophandwerk/kukanilea-git/reports/ai_sync_error.log</string>
</dict>
</plist>
